# burlap/templates/celery/celery_supervisor.template.conf
; ==============================================
;  Generated by Burlap. No not edit.
;  Celery worker Supervisor example for Django
;  Site {{ SITE }}
; ==============================================
; Based on:
; https://raw.githubusercontent.com/celery/django-celery/bf787c5d3b0de115999b876509038a416c887df1/extra/supervisord/celeryd.conf

[program:celery_{{ SITE }}]
environment=ROLE={{ ROLE }},SITE={{ SITE }}
command={{ celery_celeryd_command }}
directory={{ celery_supervisor_directory }}
user={{ celery_daemon_user }}
numprocs={{ celery_numprocs }}
; Note the string variables are for supervisor to fill in, not our templating language.
process_name=%(program_name)s_%(process_num)s
stdout_logfile={{ celery_log_path }}
stderr_logfile={{ celery_log_path }}
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998
